# alexa-airplay/rootfs/etc/s6-overlay/s6-rc.d/alexa-airplay/run
#!/command/execlineb -P
python3 -u /app/app/main.py
